# prod_deployment/bank-of-asgard-agencies-mcp.service
[Unit]
Description=Bank of Asgard Agencies MCP Server (FastMCP/SSE)
Documentation=https://github.com/fcto-demos/bank-of-asgard
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=boa
Group=boa
WorkingDirectory=/home/boa/bank-of-asgard/agencies-mcp-server
EnvironmentFile=/home/boa/bank-of-asgard/agencies-mcp-server/.env
ExecStart=/bin/bash -lc 'python server.py'

Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=bank-of-asgard-agencies-mcp

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/boa/bank-of-asgard

[Install]
WantedBy=default.target
